# micropanel.service.in
[Unit]
Description=MicroPanel OLED Menu System
After=network.target

[Service]
Type=simple
User=@SERVICE_USER@
ExecStart=/bin/sh -c '@CMAKE_INSTALL_PREFIX@/bin/micropanel @SYSTEMD_UNITFILE_ARGS@ -a -c @CMAKE_INSTALL_PREFIX@/etc/micropanel/config.json >/tmp/micropanel.log 2>&1'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target